# setup/safeskill-updater.service
[Unit]
Description=SafeSkillAgent Auto-Updater
Documentation=https://github.com/safeskill/safeskill-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/safeskill/src/setup/update.sh
Environment=SAFESKILL_INSTALL_DIR=/opt/safeskill/src
Environment=SAFESKILL_VENV_DIR=/opt/safeskill/venv
Environment=SAFESKILL_CONFIG_DIR=/etc/safeskill
Environment=SAFESKILL_LOG_DIR=/var/log/safeskill

# Run as root so it can restart the service
User=root
Group=root

# Hardening
NoNewPrivileges=yes
ProtectHome=read-only
PrivateTmp=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=safeskill-updater
