# setup/safeskill-agent.service
[Unit]
Description=SafeSkillAgent - Command Security Enforcement
Documentation=https://github.com/safeskill/safeskill-agent
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/opt/safeskill/venv/bin/safeskill start \
    --config-dir /etc/safeskill \
    --log-dir /var/log/safeskill \
    --socket /tmp/safeskill.sock
ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

User=safeskill
Group=safeskill

Environment=SAFESKILL_TRUST_MODE=normal
Environment=SAFESKILL_ENVIRONMENT=dev

WorkingDirectory=/etc/safeskill

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallArchitectures=native

ReadWritePaths=/var/log/safeskill /tmp/safeskill.sock
ReadOnlyPaths=/etc/safeskill

CapabilityBoundingSet=
AmbientCapabilities=

StandardOutput=journal
StandardError=journal
SyslogIdentifier=safeskill-agent

[Install]
WantedBy=multi-user.target
