# setup/safeskill-updater.timer
[Unit]
Description=SafeSkillAgent Auto-Update Timer
Documentation=https://github.com/safeskill/safeskill-agent

[Timer]
# Check for updates every 30 minutes
OnBootSec=5min
OnUnitActiveSec=30min

# Add random delay to avoid all machines hitting git at the same time
RandomizedDelaySec=5min

Persistent=true

[Install]
WantedBy=timers.target
